# cookbooks/postgresql/templates/default/postgresql.service.erb
[Unit]
Description=PostgreSQL

[Service]
Type=forking
ExecStart=/usr/lib/postgresql/<%= @version %>/bin/pg_ctl start -D <%= @data_directory %> -w -t <%= @pg_start_timeout %>
TimeoutStartSec=<%= @systemd_start_timeout %>
ExecStop=/usr/lib/postgresql/<%= @version %>/bin/pg_ctl stop -D <%= @data_directory %> -m fast
ExecReload=/usr/lib/postgresql/<%= @version %>/bin/pg_ctl reload -D <%= @data_directory %>
OOMScoreAdjust=-900
User=postgres
Group=postgres

[Install]
WantedBy=multi-user.target
